# keep_alive_ah.service
[Unit]
Description=Running ah on the background

[Service]
ExecStart=/bin/bash /home/ubuntu/ah-code-titans-frontend/start_ah.sh
Restart=on-failure
